# systemd/thanmatt-trash-cleanup.service
[Unit]
Description=Purge trash-cli files older than 15 days

[Service]
Type=oneshot
ExecStart=%h/thanmatt-dotfiles/scripts/trash-cleanup-linux.fish
