# systemd/thanmatt-trash-cleanup.timer
[Unit]
Description=Daily timer for trash-cleanup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
